# rootfs/etc/services.d/store/run
#!/usr/bin/with-contenv bash
###############################################
#
#  Run the custom-component-store
#
###############################################

exec communitystore